# Docker/resources/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/bin/bash -c "exec /etc/init.d/nginx start"
stderr_logfile=/var/log/supervisord/nginx-stderr.log
stdout_logfile=/var/log/supervisord/nginx-stdout.log

[program:php-fpm]
command=/bin/bash -c "exec /usr/sbin/php-fpm7.0 -F"
stderr_logfile=/var/log/supervisord/php-fpm-stderr.log
stdout_logfile=/var/log/supervisord/php-fpm-stdout.log